# addons/airsend/s6-rc.d/airsend-callback-service/run
#!/usr/bin/with-contenv bashio
cd /home
hname="$(hostname -i)"
echo "internal_url: http://${hname}:33863/"

if [ -n "${SUPERVISOR_TOKEN:-}" ]
then
	echo -n ${SUPERVISOR_TOKEN} > hass_api.token
	echo $(bashio::config 'auto_include') > auto_include.cfg
else
	echo "Not running on Home Assistant machine..."
fi

ulimit -n 4096
php -S 127.0.0.1:80 callback.php